# kernel/crontab/fitpc2_watchdog.cron
# Watchdog - Touch watchdog twice per minute
* * * * * /usr/bin/touch /dev/watchdog
* * * * * ( sleep 30 ; /usr/bin/touch /dev/watchdog )
